# cloudsnake.service
[Unit]
Description=Cloudsnake Game Server
After=network.target

[Service]
Type=simple
User=s0001311
WorkingDirectory=/home/s0001311/dev/cloudsnake
# Use the virtual environment's Python interpreter
ExecStart=/home/s0001311/dev/cloudsnake/venv/bin/python /home/s0001311/dev/cloudsnake/server.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudsnake

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
